# source/wechat-search-monitor/scripts/com.local.wechat-monitor-8765.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.local.wechat-monitor-8765</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/python3.11</string>
    <string>/Users/works14/.claude/监控/wechat-ybxhyyh-top3/run.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/works14/.claude/监控/wechat-ybxhyyh-top3</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>StandardOutPath</key>
  <string>/tmp/wechat-monitor-8765.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/wechat-monitor-8765.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>