# workbench/launchd/com.workbench.wechat-scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.workbench.wechat-scheduler</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/opt/python@3.11/libexec/bin/python3</string>
    <string>/Users/works14/Documents/zkcode/260712_SEO-GEO/workbench/scripts/wechat_scheduler.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/works14/Documents/zkcode/260712_SEO-GEO</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>WORKBENCH_BASE_URL</key>
    <string>http://127.0.0.1:8799</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>ThrottleInterval</key>
  <integer>30</integer>
  <key>StandardOutPath</key>
  <string>/Users/works14/Documents/zkcode/260712_SEO-GEO/data/runtime/wechat-scheduler.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/works14/Documents/zkcode/260712_SEO-GEO/data/runtime/wechat-scheduler.err.log</string>
</dict>
</plist>
